# templates/monasca-monitor-heartbeat.service.j2
[Unit]
Description=Monasca Monitor heartbeat
PartOf=monasca-monitor.target

[Service]
ExecStart={{ monasca_monitor_virtualenv_dir }}/bin/mm-heartbeat
Type=simple

[Install]
WantedBy=default.target
