# templates/monasca-monitor.target.j2
[Unit]
Description=Monasca Monitor service
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes

[Install]
WantedBy=multi-user.target
